# overlay/etc/profile.d/spear.sh
# Spear environment
export SPEAR=1
export SPEAR_VERSION=22.3
export SPEAR_UI_REF_W=3840
export SPEAR_UI_REF_H=2160
export SPEAR_GRAPHICS="${SPEAR_GRAPHICS:-mesa-auto}"
export AMOURANTHRTX_DEMO=0
export PATH="/usr/local/bin:${PATH}"
