# overlay/etc/profile.d/spear-field-path.sh
# Spear Field PATH — our tools win (top, nvtop, ls, …)
# C++ fieldbox + field-nvtop live in /usr/local/bin
export PATH="/usr/local/bin:${PATH}"
export SPEAR_FIELD=1
export FIELD_TOOLS=1
# Grok16 field_opt friendship
export SPEAR_CHIPS=1
